# crontab
SHELL=/bin/bash
# run python script @ 6 AM
0 6 * * * cd /app; /usr/local/bin/python /app/main.py > /proc/1/fd/1 2>&1
